# debian/sysob/DEBIAN/postinst
#!/bin/sh

mkdir -p /var/log/sysob

systemctl daemon-reload

#systemctl disable atop
#systemctl stop atop

exit 0
