# debian/sysob.prerm
#!/bin/sh

systemctl disable checkdisk.timer

