# debian/sysob/lib/systemd/system/sysob.service.d/default.conf
[Service]
Environment="SYSOB_LAST_DIR_TMP=/var/log/sysob/last_directory"
Environment="SYSOB_LOG_DIR=/var/log/sysob"
# the period of monitor
Environment="SYSOB_PERIOD=10"

 